# meta-wilhelmos/recipes-graphics/cage/files/cage-kiosk.service
[Unit]
Description=Cage Wayland kiosk running the wilhelm_renderer imgui demo
Requires=seatd.service
After=seatd.service systemd-user-sessions.service getty@tty1.service
Conflicts=getty@tty1.service
ConditionPathExists=/dev/tty1
Before=graphical.target

[Service]
User=kiosk
Group=kiosk
SupplementaryGroups=video input render seat
Environment=LIBSEAT_BACKEND=seatd
Environment=XDG_RUNTIME_DIR=/run/cage
RuntimeDirectory=cage
RuntimeDirectoryMode=0700
TTYPath=/dev/tty1
StandardInput=tty-fail
StandardOutput=journal
StandardError=journal
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
UtmpIdentifier=tty1
# seatd.service is Type=simple, so After= does not imply its socket exists
# yet; gate on actual readiness (up to 5s) to avoid a first-boot restart.
ExecStartPre=/bin/sh -c 'n=0; while [ ! -S /run/seatd.sock ] && [ $n -lt 50 ]; do sleep 0.1; n=$((n+1)); done'
ExecStart=/usr/bin/cage -ds -- /usr/bin/wilhelm-imgui-demo
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical.target
